# python-server/systemctl/homeauto-web.service
[Unit]
Description=Starts home automation webserver service
After=network-online.target

[Service]
Type=simple
ExecStart=/home/pi/home/bin/python "/home/pi/home-automation/python-server/webserver.py"

[Install]
WantedBy=multi-user.target